# api-node/supervisord.conf
[program:api]
command = node main.js environment.json
directory = /home/dotcloud/current
