# api-node/supervisord.conf
[program:api]
command = node main.js /home/dotcloud/environment.json
directory = /home/dotcloud/current
