# api-node/supervisord.conf
[program:api]
command = node main.js
directory = /home/dotcloud/current
